# systemd/pi-tv-on.timer
[Unit]
Description=Daily TV on schedule

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true
Unit=pi-tv-on.service

[Install]
WantedBy=timers.target
